# systemd/blink1.target
[Unit]
Description=Wait for blink(1) devices
