# systemd/password-agent-blink1.path
[Unit]
Description=Handle system password requests with blink(1) and USB keyfob
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target paths.target basic.target
Requires=blink1.target

[Install]
WantedBy=sysinit.target multi-user.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
